# scripts/generateStatus.timer
[Unit]
Description=Timer to generate RecRadiko's status

[Timer]
OnCalendar=*:02/10:00
Persistent=no

[Install]
WantedBy=timers.target
